# pager_duty/daemon/pager-synthetic-spooler.service
[Unit]
Description=Pager Synthetic & Sandbox Execution Spooler
After=network.target

[Service]
Type=simple
# Requires root to execute bubblewrap with new user namespaces depending on OS constraints
User=root
# Removing group specification
WorkingDirectory=/opt/odoo/addons/pager_duty/daemon
ExecStart=/usr/bin/env python3 pager_synthetic_spooler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
