# pager_duty/daemon/pager-os-monitor.service
[Unit]
Description=Generalized Pager Monitor
After=network.target

[Service]
Type=simple
User=odoo
# Removing group specification
WorkingDirectory=/opt/odoo/addons/pager_duty/daemon
EnvironmentFile=/etc/pager_daemons.env
Environment="ODOO_USER=pager_service_internal"
ExecStart=/usr/bin/env python3 generalized_monitor.py

# --- DevSecOps Systemd Sandboxing ---
# Mount the entire operating system read-only
ProtectSystem=strict
# Restrict /home to read-only so the daemon can read its config but cannot overwrite user files
ProtectHome=read-only
# Isolate /tmp and /var/tmp so malicious payloads cannot be staged
PrivateTmp=true
# Prevent access to physical hardware devices (/dev/sda, etc.)
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# Prevent the daemon from elevating privileges via setuid binaries (e.g., sudo)
NoNewPrivileges=true
# Restrict networking to basic IPv4/IPv6 and Unix sockets (blocks exotic protocols like Bluetooth or Netlink)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# Drop all root capabilities except CAP_NET_RAW, which is strictly required to execute ICMP ping sockets
CapabilityBoundingSet=CAP_NET_RAW

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
